# common/etc/crontabs/root
0 */3 * * * /etc/hosts.sh
0 1 * * * /etc/updatew.sh
