# r619ac/etc/crontabs/root
0 1 * * * /etc/updatew.sh  > /tmp/update-log.txt 2>&1
3 3 12 12 * /usr/bin/nginx-util 'check_ssl'
0 2 * * * /usr/share/mosdns/mosdns.sh geodata
